# mqtt-client-script/mqtt-client-scriptV2.service
[Unit]
Description=MQTT Client Script V2
After=network.target

[Service]
User=rt
WorkingDirectory=/home/rt/SAE23/mqtt-client-script
ExecStart=/bin/bash -c "source /home/rt/SAE23/mqtt-client-script/.venv/bin/activate && python3 /home/rt/SAE23/mqtt-client-script/mainV2.py"
Restart=always

[Install]
WantedBy=multi-user.target